# deploy/systemd/zingpop-opencode.service
[Unit]
Description=Zingpop opencode workbench backend
Documentation=https://www.zingpop.cn
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zingpop
Group=zingpop
EnvironmentFile=/etc/zingpop/zingpop.env
Environment=HOME=/srv/zingpop
Environment=XDG_DATA_HOME=/srv/zingpop/data
Environment=XDG_CONFIG_HOME=/srv/zingpop/config
Environment=XDG_CACHE_HOME=/srv/zingpop/cache
Environment=XDG_STATE_HOME=/srv/zingpop/state
WorkingDirectory=/srv/zingpop/workspaces/default
ExecStart=/bin/sh -lc 'exec "$ZINGPOP_OPENCODE_BIN" serve --hostname 127.0.0.1 --port "$ZINGPOP_PORT"'
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/srv/zingpop /var/log/zingpop

StandardOutput=append:/var/log/zingpop/opencode.log
StandardError=append:/var/log/zingpop/opencode.err.log

[Install]
WantedBy=multi-user.target
